# traits/network/tor/wifi-homepass/files/etc/uci-defaults/04_configure_wshaper
#!/bin/sh
set -euf
net=tor


uci batch <<EOF
set wshaper.settings.network=$net
set wshaper.settings.downlink=128
set wshaper.settings.uplink=128
EOF
