# traits/network/eth-dhcp/files/etc/uci-defaults/02_configure_eth_dhcp
#!/bin/sh
set -euf
itf=eth0

uci batch <<EOF
delete network.wan
set network.wan=interface
set network.wan.ifname=$itf
set network.wan.proto=dhcp
EOF
uci commit network
